# dirtytweets.cron
MAILTO=mattspitz@gmail.com
*		*		*		*		*		sudo -u mspitz python /storage/random/twittercurses/dirtytweets.py 
