# dirtytweets.cron
MAILTO=mattspitz@gmail.com
*		*		*		*		*		CONFIGFN=/storage/random/twittercurses/config.json /home/mspitz/.virtualenvs/naughtynouns/bin/python /storage/random/twittercurses/dirtytweets.py
